# spotistream.service
[Unit]
Description=Spotistream — Spotify playlist HTTP streaming server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arnonsegal
WorkingDirectory=/home/arnonsegal/scripts/spotistream
ExecStart=/home/arnonsegal/scripts/spotistream/.venv/bin/python stream.py
Restart=on-failure
RestartSec=5
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
